# src/agent_trace/autostart/linux/kimi-monitor.service.template
[Unit]
Description=AgentTrace - AI IDE Session Monitor for Coze Compass
Documentation=https://github.com/your-repo/agent-trace
After=network.target
Wants=network.target

[Service]
Type=simple
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{WORKING_DIR}}

# Load environment variables from secure file
# Note: Create this file with your API credentials
EnvironmentFile={{ENV_FILE_PATH}}

# Non-sensitive environment variables
Environment="PYTHONPATH={{PYTHONPATH}}"

# Python executable and module
ExecStart={{PYTHON_PATH}} -m agent_trace --log-level INFO
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security
NoNewPrivileges=false
ProtectSystem=no
ProtectHome=no

# Logging
StandardOutput=append:{{LOG_DIR}}/agent-trace.stdout.log
StandardError=append:{{LOG_DIR}}/agent-trace.stderr.log
SyslogIdentifier=agent-trace

[Install]
WantedBy=multi-user.target
